# lesson_10/template/watchlog.timer.j2
[Unit]
Description=Run watchlog script every 30 second

[Timer]
# Run every 30 second
OnBootSec=1min
OnUnitActiveSec=30
Unit=watchlog.service
[Install]
WantedBy=timers.target
